# deploy/systemd/executor.service
[Unit]
Description=Distributed Prediction Market - executor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=market
Group=market
WorkingDirectory=/home/market/Distributed-Prediction-Market-Server
EnvironmentFile=-/etc/prediction-market.env
ExecStart=/home/market/Distributed-Prediction-Market-Server/executor/run_forever.sh
Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=10
StandardOutput=append:/var/log/prediction-market/executor.log
StandardError=append:/var/log/prediction-market/executor.log

[Install]
WantedBy=multi-user.target
